# etc/profile.d/qgis.sh
export PYTHONPATH=/opt/qgis/share/qgis/python
